# src/infra/com.aura.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aura.telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/oxyzen/.local/bin/claude-telegram-bot</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/oxyzen/claude-code-telegram</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/oxyzen/claude-code-telegram/logs/bot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/oxyzen/claude-code-telegram/logs/bot.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/oxyzen/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/oxyzen</string>
        <key>ORT_DISABLE_COREML</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
