# src/infra/com.aura.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aura.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/oxyzen/claude-code-telegram/bin/aura</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/var/log/aura/bot.out</string>
    <key>StandardErrorPath</key>
    <string>/var/log/aura/bot.err</string>
</dict>
</plist>
