# menubar/com.aura.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aura.menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/oxyzen/.local/share/uv/tools/claude-code-telegram/bin/python3</string>
        <string>/Users/oxyzen/claude-code-telegram/menubar/aura_bar.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/aura-menubar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aura-menubar-err.log</string>
</dict>
</plist>
